# systemd/pathogen-update.service
[Unit]
Description=Update Vim Pathogen Plugins

[Service]
ExecStart=/home/jabrouwer/pathogen-update.sh
