# systemd/trash-empty.service
[Unit]
Description=Empty Old Trash

[Service]
ExecStart=/home/jabrouwer/empty-trash.sh
