# systemd/trash-empty.service
[Unit]
Description=Empty Month Old Trash

[Service]
Type=oneshot
ExecStart=/home/jabrouwer/trash-empty.sh
